# src/rootfs/etc/services.d/davmail/run
#!/bin/sh
set -e -u

exec s6-setuidgid nobody /opt/davmail/davmail.sh /etc/davmail/davmail.properties
