# src/rootfs/etc/services.d/davmail/run
#!/bin/sh

/usr/lib/davmail/davmail.sh /etc/davmail/davmail.properties
