# src/rootfs/etc/services.d/davmail/run
#!/bin/sh
set -e

/usr/share/davmail/davmail.sh /etc/davmail/davmail.properties
